# root/etc/services.d/vaultwarden/run
#!/usr/bin/with-contenv bash

cd / || exit

echo "Starting Vaultwarden."
redirect_cmd exec s6-setuidgid abc \
	vaultwarden
